# nercone-webserver-autoupdater.timer
[Unit]
Description=Nercone Web Server Auto-updater

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
